# initscripts/cloud
 ### BEGIN INIT INFO
 # Provides:          cloud 
 # Required-Start:    $remote_fs $syslog
 # Required-Stop:     $remote_fs $syslog
 # Default-Start:     2 3 4 5
 # Default-Stop:      0 1 6
 # Short-Description: Cloud initscript
 # Description:       This file should be used to construct scripts to be
 #                    placed in /etc/init.d.
 ### END INIT INFO
  
case "$1" in
    start)
        echo "Starting Cloud..."
        # run application you want to start
        # python3 -u `which cloud` run >> /var/log/cloud.log &
        cloud run >> /var/log/cloud.log &
        ;;
    stop)
        echo "Stopping Cloud..."
        # kill application you want to stop
        killall cloud
        ;;
    *)
        echo "Usage: /etc/init.d/cloud{start|stop}"
        exit 1
        ;;
esac

exit 0
